# Pricing_Pipelling/.ipynb_checkpoints/old/cronjob/uber_job.cron
* * * * * /home/admin1/anaconda/bin/python /home/admin1/Galvanize/FinalProjectDE/cronjob/uber_job_scheduler.py > /home/admin1/uber.log

